# deploy/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/home/Synth/deploy/gunicorn.sock

[Install]
WantedBy=socket.target

