# deploy/think.gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/home/Synth/deploy/think.gunicorn.socket

[Install]
WantedBy=socket.target

